# references/nemoclaw/overlay/usr/lib/systemd/system/nemoclaw-gateway.service
[Unit]
Description=NemoClaw AI Agent Gateway
After=nemoclaw-setup.service docker.service network-online.target
Requires=docker.service nemoclaw-setup.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/nemoclaw/gateway.env
ExecStartPre=/usr/bin/docker rm -f nemoclaw-gateway
ExecStart=/usr/lib/nemoclaw/nemoclaw-gateway-start
ExecStop=/usr/bin/docker stop nemoclaw-gateway
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
